# packages/ccproxy-api/ccproxy_api-0.1.7.tar.gz/ccproxy_api-0.1.7/systemd/ccproxy.service.template
[Unit]
Description=CCProxy API Server
After=network.target

[Service]
Type=simple
WorkingDirectory={{WORKING_DIR}}
ExecStart={{UV_PATH}} run ccproxy serve
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ccproxy
Environment="PATH={{USER_PATH}}"
Environment="HOME={{USER_HOME}}"
Environment="UV_PROJECT={{UV_PROJECT}}"
{{EXTRA_ENV}}

[Install]
WantedBy=default.target
