# packages/automagik-hive/automagik_hive-0.1.0a59-py3-none-any.whl/docker/agent/supervisord.conf
[supervisord]
nodaemon=true
user=agent
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700
chown=agent:agent

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# PostgreSQL Service
[program:postgresql]
command=/usr/lib/postgresql/16/bin/postgres -D /var/lib/postgresql/data/pgdata -c config_file=/etc/postgresql/16/main/postgresql.conf
user=agent
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/postgresql/postgresql.log
stderr_logfile=/var/log/postgresql/postgresql.error.log
environment=PGDATA="/var/lib/postgresql/data/pgdata"
priority=100

# FastAPI Application Service
[program:hive-agent]
command=python -m uvicorn api.serve:app --host 0.0.0.0 --port 38886 --workers 2
directory=/app
user=agent
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/hive/agent-api.log
stderr_logfile=/var/log/hive/agent-api.error.log
environment=HIVE_DATABASE_URL="postgresql+psycopg://agent:agent_secure_password@localhost:5432/hive_agent",HIVE_API_PORT="38886",HIVE_SERVICE_MODE="agent"
priority=200

# Health Check Service (Optional)
[program:health-monitor]
command=/app/docker/agent/health-monitor.sh
user=agent
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/hive/health-monitor.log
priority=300