# packages/genesis-notification/genesis_notification-0.3.0.tar.gz/genesis_notification-0.3.0/etc/systemd/genesis-notification-builder-agent.service
[Unit]
Description=Genesis Notification Builder Agent
After=network-online.target

[Service]
TimeoutStopSec=5
Restart=on-failure
RestartSec=5s
KillSignal=SIGINT
ExecStartPre=/bin/sleep 2
ExecStart=genesis-notification-builder-agent --config-file /etc/genesis_notification/genesis_notification.conf

[Install]
WantedBy=multi-user.target
